# policyd.cron
#!/bin/sh

exec /etc/rc.d/init.d/policyd cron
